# beautyai-webui.service
[Unit]
Description=BeautyAI Web UI Service (Frontend)
After=network.target beautyai-api.service
Wants=network.target
Requires=beautyai-api.service

[Service]
Type=simple
User=lumi
Group=lumi
WorkingDirectory=/home/lumi/beautyai/frontend/src

# Load environment variables from production file
EnvironmentFile=-/home/lumi/beautyai/.env.production

# Core environment settings
Environment=PYTHONPATH=/home/lumi/beautyai/frontend/src:/home/lumi/beautyai/frontend/venv/lib/python3.12/site-packages
Environment=PYTHONUNBUFFERED=1
Environment=VIRTUAL_ENV=/home/lumi/beautyai/frontend/venv
Environment=FLASK_ENV=production

# BeautyAI specific configuration - Production mode
Environment=PORT=5001
Environment=BEHIND_PROXY=true

ExecStart=/home/lumi/beautyai/frontend/venv/bin/python /home/lumi/beautyai/frontend/src/app.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=beautyai-webui

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/log

[Install]
WantedBy=multi-user.target
