# beautyai-api.service
[Unit]
Description=BeautyAI Inference API Server (Development)
After=network.target
Wants=network.target

[Service]
Type=simple
User=lumi
Group=lumi
WorkingDirectory=/home/lumi/beautyai
Environment=PATH=/home/lumi/beautyai/venv/bin
Environment=PYTHONPATH=/home/lumi/beautyai
Environment=CUDA_VISIBLE_DEVICES=0
ExecStart=/home/lumi/beautyai/venv/bin/uvicorn beautyai_inference.api.app:app --host 0.0.0.0 --port 8000 --reload
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGINT
TimeoutStopSec=5
Restart=on-failure
RestartSec=3

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/home/lumi/beautyai
ProtectHome=read-only

# Resource limits
LimitNOFILE=65536
MemoryMax=32G

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=beautyai-api

[Install]
WantedBy=multi-user.target
