# beautyai-api.service
[Unit]
Description=BeautyAI Inference API Server (Development)
After=network.target
Wants=network.target

[Service]
Type=simple
User=lumi
Group=lumi
WorkingDirectory=/home/lumi/beautyai/backend
Environment=PATH=/home/lumi/beautyai/backend/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=PYTHONPATH=/home/lumi/beautyai/backend
Environment=CUDA_VISIBLE_DEVICES=0
Environment=VOICE_STREAMING_ENABLED=1
Environment=VOICE_STREAMING_PHASE4=1
ExecStart=/home/lumi/beautyai/backend/venv/bin/python /home/lumi/beautyai/backend/run_server.py
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGINT
TimeoutStopSec=5
Restart=on-failure
RestartSec=3

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/home/lumi/beautyai/backend
ProtectHome=read-only

# Resource limits
LimitNOFILE=65536
MemoryMax=32G

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=beautyai-api

[Install]
WantedBy=multi-user.target
